# ansible/roles/create_proxmox_vm/files/rc.local
#!/bin/bash
dpkg-reconfigure openssh-server
sed -e 's/\s*\([\+0-9a-zA-Z]*\).*/\1/' << EOF | fdisk /dev/vda
  d #
  1
  n # new partition
  p # primary partition
  1 # partition number 1
    # default - start at beginning of disk
    # 100 MB boot parttion
  w # write the partition table
  q # and we're done
EOF
partprobe /dev/vda
resize2fs /dev/vda1
>/etc/rc.local

